# build/package/freebsd/relique-client/pkg-plist
@owner relique
@group relique
@mode 0644

@(,,755) bin/relique-client
@(,,755) etc/rc.d/relique-client

etc/relique/certs/cert.pem
etc/relique/certs/key.pem

@sample etc/relique/client.toml.sample
@dir(,,0755) etc/relique
@dir(,,0755) etc/relique/certs

@dir(,,0755) relique/modules
@dir(,,0755) relique/modules/generic
relique/modules/generic/README.md
relique/modules/generic/default.toml
