# build/package/freebsd/relique-server/pkg-plist
@owner relique
@group relique
@mode 0644

@dir(,,0755) etc/relique
@dir(,,0755) etc/relique/certs
@dir(,,0755) etc/relique/clients
@dir(,,0755) etc/relique/schedules
@dir(,,0755) relique/db
@dir(,,0755) relique/storage
@dir(,,0755) relique/modules
@dir(,,0755) relique/modules/generic

@(,,755) bin/relique-server
@(,,755) etc/rc.d/relique-server

etc/relique/certs/cert.pem
etc/relique/certs/key.pem
@sample etc/relique/server.toml.sample
etc/relique/clients/example.toml.disabled
etc/relique/schedules/daily.toml
etc/relique/schedules/manual.toml
etc/relique/schedules/weekly.toml

relique/modules/generic/README.md
relique/modules/generic/default.toml
@comment TODO: Convert to sample to avoid overwrite on update
