# build/init/relique-server.service
[Unit]
Description=Relique server daemon
After=network.target

[Service]
Type=simple
User=relique
Group=relique
ExecStart=/usr/bin/relique-server start
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
